# public/autostart/.config/autostart/CEB36C9B3BF99F11BA499A34F5BC7C651BB39D2A.chromium-service.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/usr/bin/false
Name=Chromium Service
