# public/autostart/.config/autostart/startup.sh.desktop
[Desktop Entry]
Type=Application
Name=startup.sh
Exec=/home/andrew/Scripts/startup.sh
